# script/smith.service
[Unit]
Description=smith

[Service]
Type=forking
PIDFile=/var/run/smith.pid
ExecStart=/opt/smith/smithworker1
Restart=always

[Install]
WantedBy=multi-user.target